# experiment.service
[Unit]
Description=File Server service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=python3 -u /root/experiment/files/experiment.py -l local -i IPGOESHERE
StandardOutput=file:/tmp/experiment.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
